# git/deploy/workersproxy.service
[Unit]
Description=cdn-git-proxy (raw.githubusercontent.com reverse proxy)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/cdn-git-proxy
EnvironmentFile=/opt/cdn-git-proxy/.env
ExecStart=/usr/bin/node /opt/cdn-git-proxy/server.js
Restart=on-failure
RestartSec=2

# Run as a dedicated unprivileged user if you want (recommended):
# User=cdn-git-proxy
# Group=cdn-git-proxy

[Install]
WantedBy=multi-user.target

